# Linux/bluetooth-monitor.timer
# Install to: /etc/systemd/system/bluetooth-monitor.timer
# Then:
#   sudo systemctl daemon-reload
#   sudo systemctl enable --now bluetooth-monitor.timer

[Unit]
Description=Run Bluetooth Monitor Script Every Minute

[Timer]
OnBootSec=1min
OnUnitActiveSec=1min
Unit=bluetooth-monitor.service

[Install]
WantedBy=timers.target
